# install/systemd/wi-lab.service.template
[Unit]
Description=Wi-Lab WiFi Access Point Manager
Documentation=https://github.com/caporalesimone/wi-lab
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=WILAB_DIR
ExecStart=VENV_PATH/bin/python WILAB_DIR/main.py
Restart=always
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s
Environment="CONFIG_PATH=WILAB_DIR/config.yaml"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wilab

[Install]
WantedBy=multi-user.target
